# supervisor.conf
[supervisord]
nodaemon=true

[program:app1]
directory=/home/client
command=npm start
environment=PORT="3000"
stdout_logfile=/logs/client1.log
stdout_logfile_maxbytes=1GB
redirect_stderr=true